# drf_back/config/django.conf
[program:django]
command = 
directory = 
user = 
autorestart=true
redirect_stderr = 
